# systemd/django.socket
[Unit]
Description=Example web app socket

[Socket]
ListenStream=/run/django.sock

[Install]
WantedBy=sockets.target
